# scripts/com.matt.samwise-2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.matt.samwise-2</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mjohnst/samwise/Personal-Apps/samwise-2/scripts/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mjohnst/samwise/Personal-Apps/samwise-2</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mjohnst/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/mjohnst</string>
        <key>PORT</key>
        <string>8090</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/Users/mjohnst/Library/Logs/samwise-2.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mjohnst/Library/Logs/samwise-2.err.log</string>
</dict>
</plist>
